# scripts/supervisor/job-worker.ini
[program:all_except_invest]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-job-worker-all/bin/ttsd-job-worker

[program:loan_mq_consumer]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-loan-mq-consumer/bin/ttsd-loan-mq-consumer
stopwaitsecs=60

[program:message_mq_consumer]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-message-mq-consumer/bin/ttsd-message-mq-consumer
stopwaitsecs=60

[program:point_mq_consumer]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-point-mq-consumer/bin/ttsd-point-mq-consumer
stopwaitsecs=60

[program:user_mq_consumer]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-user-mq-consumer/bin/ttsd-user-mq-consumer
stopwaitsecs=60

[program:activity_mq_consumer]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-activity-mq-consumer/bin/ttsd-activity-mq-consumer
stopwaitsecs=60

[program:auditLog_mq_consumer]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-auditLog-mq-consumer/bin/ttsd-auditLog-mq-consumer
stopwaitsecs=60

[program:amount_mq_consumer]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-amount-mq-consumer/bin/ttsd-amount-mq-consumer
stopwaitsecs=60

[program:eamil_mq_consumer]
environment=TTSD_ETCD_ENV=prod
command=/workspace/ttsd-email-mq-consumer/bin/ttsd-email-mq-consumer
stopwaitsecs=60