# scripts/supervisor/job-worker.ini
[program:all_except_invest]
command=/workspace/ttsd-job-worker-all/bin/ttsd-job-worker

[program:jpush]
command=/workspace/ttsd-job-worker-jpush/bin/ttsd-job-worker

[program:repay]
command=/workspace/ttsd-job-worker-repay/bin/ttsd-job-worker

[program:loan_mq_consumer]
command=/workspace/ttsd-loan-mq-consumer/bin/ttsd-loan-mq-consumer
stopsignal=TERM
stopwaitsecs=60

[program:message_mq_consumer]
command=/workspace/ttsd-message-mq-consumer/bin/ttsd-message-mq-consumer
stopsignal=TERM
stopwaitsecs=60

[program:point_mq_consumer]
command=/workspace/ttsd-point-mq-consumer/bin/ttsd-point-mq-consumer
stopsignal=TERM
stopwaitsecs=60

[program:user_mq_consumer]
command=/workspace/ttsd-user-mq-consumer/bin/ttsd-user-mq-consumer

[program:diagnosis]
command=/workspace/ttsd-diagnosis/bin/ttsd-diagnosis
