# scripts/supervisor/job-worker.ini
[program:all_except_invest]
command=/workspace/ttsd-job-worker-all/bin/ttsd-job-worker

[program:invest]
command=/workspace/ttsd-job-worker-invest/bin/ttsd-job-worker

[program:jpush]
command=/workspace/ttsd-job-worker-jpush/bin/ttsd-job-worker

[program:loan_mq_consumer]
command=/workspace/ttsd-loan-mq-consumer/bin/ttsd-loan-mq-consumer
